# konker.samples/rest-json-enrichment/product-data.supervisord.conf
[program:product-data]
command=/path/to/virtualenv/bin/gunicorn -b ip:port product-data:app
directory=/path/to/product-data-dir
user=nobody
autostart=true
autorestart=true
redirect_stderr=true
